# debian_overlay/etc/rc.local
#!/bin/bash

if [ -f /usr/local/sbin/start.sh ]; then
	/usr/local/sbin/start.sh &
fi

exit 0

